# exspider/supervisor_conf/storekeeper/storekeeper_prod.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=/var/run/supervisor_storekeeper.sock   ; the path to the socket file
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=root              ; default is no username (open server)
;password=root               ; default is no password (open server)

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
logfile=/tmp/storekeeper.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=1           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/storekeeper.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
;umask=022                   ; process file creation umask; default 022
;user=chrism                 ; default is current user, required if root
;identifier=supervisor       ; supervisord identifier, default is 'supervisor'
;directory=/tmp              ; default is not to cd during start
;nocleanup=true              ; don't clean up tempfiles at start; default false
;childlogdir=/tmp            ; 'AUTO' child log dir, default $TEMP
;environment=KEY="value"     ; key value pairs to add to environment
;strip_ansi=false            ; strip ansi escape codes in logs; def. false

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix:///var/run/supervisor_storekeeper.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The sample program section below shows all possible program subsection values.
; Create one or more 'real' program: sections to be able to control them under
; supervisor.

;[program:theprogramname]
;command=/bin/cat              ; the program (relative uses PATH, can take args)
;process_name=%(program_name)s ; process_name expr (default %(program_name)s)
;numprocs=1                    ; number of processes copies to start (def 1)
;directory=/tmp                ; directory to cwd to before exec (def no cwd)
;umask=022                     ; umask for process (default None)
;priority=999                  ; the relative start priority (default 999)
;autostart=true                ; start at supervisord start (default: true)
;startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
;startretries=3                ; max # of serial start failures when starting (default 3)
;autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
;exitcodes=0,2                 ; 'expected' exit codes used with autorestart (default 0,2)
;stopsignal=QUIT               ; signal used to kill process (default TERM)
;stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
;stopasgroup=false             ; send stop signal to the UNIX process group (default false)
;killasgroup=false             ; SIGKILL the UNIX process group (def false)
;user=chrism                   ; setuid to this UNIX account to run the program
;redirect_stderr=true          ; redirect proc stderr to stdout (default false)
;stdout_logfile=/a/path        ; stdout log path, NONE for none; default AUTO
;stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
;stdout_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stdout_events_enabled=false   ; emit events on stdout writes (default false)
;stderr_logfile=/a/path        ; stderr log path, NONE for none; default AUTO
;stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
;stderr_logfile_backups=10     ; # of stderr logfile backups (0 means none, default 10)
;stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
;stderr_events_enabled=false   ; emit events on stderr writes (default false)
;environment=A="1",B="2"       ; process environment additions (def no adds)
;serverurl=AUTO                ; override serverurl computation (childutils)

; The sample eventlistener section below shows all possible eventlistener
; subsection values.  Create one or more 'real' eventlistener: sections to be
; able to handle event notifications sent by supervisord.


;;;;;;;;;;;;;;;;;;        prometheus       ;;;;;;;;;;;;;;;;;;;;;;;;;;
;[program:care_server]
;command=python manage.py start_prometheus_server
;directory=/opt/exspider/
;autorestart=true
;stopsignal=QUIT
;process_name=top
;redirect_stderr=true
;redirect_stdout=true
;stopasgroup=true  ;杀掉整个group，防止出现孤儿进程
;killasgroup=true  ;杀掉整个group，防止出现孤儿进程
;stdout_logfile=/opt/logs/care_server.log
;stderr_logfile=/opt/logs/care_server.log


[program:storekeeper_binance]
command=python manage.py store_keeper prod_20001 spider2 market.kline.binance.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_binance.log
stderr_logfile=/opt/logs/storekeeper_binance.error.log


[program:storekeeper_huobipro]
command=python manage.py store_keeper prod_20002 spider2 market.kline.huobipro.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_huobipro.log
stderr_logfile=/opt/logs/storekeeper_huobipro.error.log

[program:storekeeper_okex]
command=python manage.py store_keeper prod_20003 spider2 market.kline.okex.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_okex.log
stderr_logfile=/opt/logs/storekeeper_okex_error.log

[program:storekeeper_bitfinex]
command=python manage.py store_keeper prod_20004 spider2 market.kline.bitfinex.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bitfinex.log
stderr_logfile=/opt/logs/storekeeper_bitfinex_error.log

[program:storekeeper_bitmex]
command=python manage.py store_keeper prod_20005 spider2 market.kline.bitmex.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bitmex.log
stderr_logfile=/opt/logs/storekeeper_bitmex.error.log


[program:storekeeper_hotbit]
command=python manage.py store_keeper prod_20006 spider1 market.kline.hotbit.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_hotbit.log
stderr_logfile=/opt/logs/storekeeper_hotbit.error.log


[program:storekeeper_gateio]
command=python manage.py store_keeper prod_20007 spider1 market.kline.gateio.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_gateio.log
stderr_logfile=/opt/logs/storekeeper_gateio.error.log


[program:storekeeper_bibox]
command=python manage.py store_keeper prod_20008 spider1 market.kline.bibox.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bibox.log
stderr_logfile=/opt/logs/storekeeper_bibox.error.log

[program:storekeeper_coinbasepro]
command=python manage.py store_keeper prod_20011 spider1 market.trade.coinbasepro.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_coinbasepro.log
stderr_logfile=/opt/logs/storekeeper_coinbasepro.error.log

[program:storekeeper_bittrex]
command=python manage.py store_keeper prod_20012 spider1 market.trade.bittrex.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bittrex.log
stderr_logfile=/opt/logs/storekeeper_bittrex.error.log

[program:storekeeper_hitbtc]
command=python manage.py store_keeper prod_20009 spider3 market.kline.hitbtc.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_hitbtc.log
stderr_logfile=/opt/logs/storekeeper_hitbtc.error.log

[program:storekeeper_kraken]
command=python manage.py store_keeper prod_20010 spider3 market.kline.kraken.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_kraken.log
stderr_logfile=/opt/logs/storekeeper_kraken.error.log

[program:storekeeper_kucoin]
command=python manage.py store_keeper prod_20013 spider3 market.trade.kucoin.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_kucoin.log
stderr_logfile=/opt/logs/storekeeper_kucoin.error.log

[program:storekeeper_bitstamp]
command=python manage.py store_keeper prod_20014 spider3 market.trade.bitstamp.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bitstamp.log
stderr_logfile=/opt/logs/storekeeper_bitstamp.error.log

[program:storekeeper_bithumb]
command=python manage.py store_keeper prod_20015 spider3 market.trade.bithumb.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=5MB
stdout_logfile=/opt/logs/storekeeper_bithumb.log
stderr_logfile=/opt/logs/storekeeper_bithumb.error.log

[program:storekeeper_mxc]
command=python manage.py store_keeper prod_20016 spider1 market.kline.mxc.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=1MB
stdout_logfile=/opt/logs/storekeeper_mxc.log
stderr_logfile=/opt/logs/storekeeper_mxc.error.log

[program:storekeeper_idax]
command=python manage.py store_keeper prod_20017 spider3 market.kline.idax.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=1MB
stdout_logfile=/opt/logs/storekeeper_idax.log
stderr_logfile=/opt/logs/storekeeper_idax.error.log


[program:storekeeper_biki]
command=python manage.py store_keeper prod_20018 spider3 market.kline.biki.*.*.*.*.*  storekeeper_prod
directory=/opt/exspider/
redirect_stderr=true
redirect_stdout=true
stdout_logfile_maxbytes=1MB
stdout_logfile=/opt/logs/storekeeper_biki.log
stderr_logfile=/opt/logs/storekeeper_biki.error.log